# autostart/chrome-famdcdojlmjefmhdpbpmekhodagkodei-Default.desktop
#!/usr/bin/env xdg-open
[Desktop Entry]
Version=1.0
Terminal=false
Type=Application
Name=Outlook (PWA)
MimeType=x-scheme-handler/mailto
Exec=/opt/google/chrome/google-chrome --profile-directory=Default --app-id=famdcdojlmjefmhdpbpmekhodagkodei %U
Icon=chrome-famdcdojlmjefmhdpbpmekhodagkodei-Default
StartupWMClass=crx_famdcdojlmjefmhdpbpmekhodagkodei
Actions=New-event;New-message;Open-calendar

[Desktop Action New-event]
Name=New event
Exec=/opt/google/chrome/google-chrome --profile-directory=Default --app-id=famdcdojlmjefmhdpbpmekhodagkodei --app-launch-url-for-shortcuts-menu-item=https://outlook.office365.us/calendar/deeplink/compose

[Desktop Action New-message]
Name=New message
Exec=/opt/google/chrome/google-chrome --profile-directory=Default --app-id=famdcdojlmjefmhdpbpmekhodagkodei --app-launch-url-for-shortcuts-menu-item=https://outlook.office365.us/mail/deeplink/compose

[Desktop Action Open-calendar]
Name=Open calendar
Exec=/opt/google/chrome/google-chrome --profile-directory=Default --app-id=famdcdojlmjefmhdpbpmekhodagkodei --app-launch-url-for-shortcuts-menu-item=https://outlook.office365.us/calendar
